# scheduled_task/crontab
#run python script every minutes
*/20 * * * * python3 /app/hello.py > /proc/1/fd/1 2>/proc/1/fd/2
